# server_config/crontab
* * * * * /srv/pypln/project/deploy/exec_in_virtualenv.sh /srv/pypln/project/web/pypln/web/manage.py update_index --settings=settings.production
